# files/clear-docker-cache.cron
#!/bin/sh

[ -x /usr/share/gitlab-runner/clear-docker-cache ] && /usr/share/gitlab-runner/clear-docker-cache >/dev/null
